# raumfeld/rootfs-base-armada/lib/udev/rules.d/99-wlan0.rules
KERNEL=="wlan0", RUN+="/usr/sbin/iw dev wlan0 set txpower fixed 1800"
