# raumfeld/rootfs-geode/etc/udev/rules.d/20-usbunmount.rules
KERNEL=="sd??", ACTION=="remove", RUN+="/raumfeld/meta-server/handleUnmount.sh"
