# raumfeld/rootfs-geode/etc/udev/rules.d/80-usbmount.rules
KERNEL=="sd??", ACTION=="add", RUN+="/raumfeld/meta-server/handleMount.sh"

