# raumfeld/rootfs-geode/etc/udev/rules.d/80-usbmount.rules
KERNEL=="sda?", BUS=="usb", ACTION=="add",    RUN+="/sbin/usbmount add"
KERNEL=="sda?",             ACTION=="remove", RUN+="/sbin/usbmount remove"
