# domestic_recorder/fs-overlay/etc/udev/rules.d/100-wpa_supplicant.rules
SUBSYSTEM=="net",ACTION=="add",KERNEL=="wlan0",RUN+="/bin/systemctl start wpa_supplicant"
SUBSYSTEM=="net",ACTION=="remove",KERNEL=="wpan0",RUN+="/bin/systemctl stop wpa_supplicant"
