# raumfeld/rootfs-remotecontrol-arm/etc/udev/rules.d/20-usbunmount.rules
KERNEL=="sd??", ACTION=="remove", RUN+="/usr/share/remote-control/scripts/unmount.sh"
