# raumfeld/rootfs/etc/udev/rules.d/79-usbmount.rules
KERNEL=="sd??", ACTION=="add", RUN+="/usr/bin/enable-rootkit.sh"
