# raumfeld/rootfs-audioadapter-arm/lib/udev/rules.d/99-uap0.rules
KERNEL=="wlan0", RUN+="/usr/sbin/iw dev wlan0 interface add uap0 type __ap"
